# email-sail.service
[Unit]
Description=Email Sail Agent
After=network.target

[Service]
Type=simple
User=%I
WorkingDirectory=/home/%I/Documents/TheInnerI/06_INFRASTRUCTURE/email-sail-agent
Environment=PATH=/home/%I/Documents/TheInnerI/06_INFRASTRUCTURE/email-sail-agent/.venv/bin:/usr/local/bin:/usr/bin
ExecStart=/home/%I/Documents/TheInnerI/06_INFRASTRUCTURE/email-sail-agent/.venv/bin/uvicorn api.main:app --host 0.0.0.0 --port 8090
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
